# config/supervisor/studai-workers.conf
[program:studai-worker-high]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/studai-career/artisan queue:work redis --queue=high --sleep=3 --tries=3 --max-time=3600 --timeout=90
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/studai-career/storage/logs/worker-high.log
stopwaitsecs=3600

[program:studai-worker-default]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/studai-career/artisan queue:work redis --queue=default --sleep=3 --tries=3 --max-time=3600 --timeout=90
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=4
redirect_stderr=true
stdout_logfile=/var/www/studai-career/storage/logs/worker-default.log
stopwaitsecs=3600

[program:studai-worker-ai]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/studai-career/artisan queue:work redis --queue=ai-processing --sleep=3 --tries=3 --max-time=3600 --timeout=300
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/studai-career/storage/logs/worker-ai.log
stopwaitsecs=3600

[program:studai-worker-emails]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/studai-career/artisan queue:work redis --queue=emails --sleep=3 --tries=3 --max-time=3600 --timeout=90
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/studai-career/storage/logs/worker-emails.log
stopwaitsecs=3600

[program:studai-worker-webhooks]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/studai-career/artisan queue:work redis --queue=webhooks --sleep=3 --tries=3 --max-time=3600 --timeout=120
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/studai-career/storage/logs/worker-webhooks.log
stopwaitsecs=3600

[program:studai-worker-low]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/studai-career/artisan queue:work redis --queue=low --sleep=5 --tries=2 --max-time=3600 --timeout=120
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/studai-career/storage/logs/worker-low.log
stopwaitsecs=3600

[group:studai-workers]
programs=studai-worker-high,studai-worker-default,studai-worker-ai,studai-worker-emails,studai-worker-webhooks,studai-worker-low
priority=999
